# supervisord.conf
[supervisord]
nodaemon=true

[program:startioc]
command=softIoc -d testioc.db
priority=1
